# firmware/files/etc/crontabs/root
0-59 * * * * /save_time.sh
0,30 * * * * /generate_post.sh
30 3 * * * * reboot
